# devops/grade.cron
5 22 * * * oasis /home/oasis/hwcentral/grader/scripts/grade-overnight.sh