# devops/sleep.cron
0 22 * * * root  /home/oasis/hwcentral/devops/sleep-mode-on.sh
0 4  * * * root  /home/oasis/hwcentral/devops/sleep-mode-off.sh
