# devops/crontabs/notify.cron
0 17 * * * oasis /home/oasis/hwcentral/pylon/scripts/notify-due.sh >> /home/oasis/hwcentral/devops/cron.log 2>&1

