# devops/crontabs/sleep.cron
0 22 * * * root  /home/oasis/hwcentral/devops/sleep-mode-on.sh  >> /home/oasis/hwcentral/devops/cron.log
0 4  * * * root  /home/oasis/hwcentral/devops/sleep-mode-off.sh >> /home/oasis/hwcentral/devops/cron.log

